# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/data/supervisord.log
logfile_maxbytes=5MB
logfile_backups=2
pidfile=/tmp/supervisord.pid

[program:engine]
command=python3 /app/main.py
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=15
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/app/data/engine.log
stderr_logfile=/app/data/engine_err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

[program:api]
command=uvicorn dashboard.api:app --host 0.0.0.0 --port 8000 --log-level info
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/app/data/api.log
stderr_logfile=/app/data/api_err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
