# nzt48-aegis-v2/crontab
# RT2: Bridge Health Monitor — every 15 min Mon-Fri during market hours (Telegram alert if unhealthy)
*/15 8-20 * * 1-5 cd /app && python3 -m python_brain.ouroboros.bridge_health >> /var/log/bridge_health.log 2>&1

# N10b: External Monitor — every 5 min Mon-Fri (deep health check: engine, disk, IBKR, Redis, cron)
*/5 * * * 1-5 cd /app && python3 -m python_brain.ouroboros.external_monitor --quiet >> /var/log/external_monitor.log 2>&1

# N10b: Daily Health Report — 08:00 UTC Mon-Fri (Telegram summary)
0 8 * * 1-5 cd /app && python3 -m python_brain.ouroboros.external_monitor --daily >> /var/log/external_monitor.log 2>&1

# N10c: Log Rotation — 04:45 UTC daily (before nightly learning, keeps 7 days)
45 4 * * * cd /app && python3 -m python_brain.ouroboros.log_rotate >> /var/log/log_rotate.log 2>&1

# Ouroboros v6.0 Nightly Learning Loop — 23:50 ET = 04:50 UTC (EST) Mon-Fri
# MUST run during DARK window (21:00-23:00 London = 21:00-23:00 UTC GMT / 22:00-00:00 BST)
50 4 * * 1-5 cd /app && python3 -m python_brain.ouroboros.nightly_v6 >> /var/log/ouroboros.log 2>&1

# Ouroboros Config Writer — 04:51 UTC Mon-Fri (right after nightly_v6, writes dynamic_weights.toml)
51 4 * * 1-5 cd /app && python3 -m python_brain.ouroboros.config_writer >> /var/log/config_writer.log 2>&1

# N4b: Win/Loss Indicator Delta — 04:52 UTC Mon-Fri (after nightly_v6, feeds Win_Loss_Delta tab)
52 4 * * 1-5 cd /app && python3 -m python_brain.ouroboros.win_loss_delta --push-sheets >> /var/log/win_loss_delta.log 2>&1

# N6a: Claude Nightly Review — 04:53 UTC Mon-Fri (after nightly_v6 + config_writer + delta, read-only)
53 4 * * 1-5 cd /app && python3 -m python_brain.ouroboros.claude_review --send-telegram >> /var/log/claude_review.log 2>&1

# N6b: Claude Morning Briefing — 07:45 UTC Mon-Fri (before European session, uses N6a output)
45 7 * * 1-5 cd /app && python3 -m python_brain.ouroboros.claude_briefing --send-telegram >> /var/log/claude_briefing.log 2>&1

# Ouroboros v6.0 Backfill Simulator — 07:00 UTC daily Mon-Fri (pre-market learning)
0 7 * * 1-5 cd /app && python3 -m python_brain.ouroboros.backfill_simulator >> /var/log/ouroboros_backfill.log 2>&1

# ISA Universe Refresh — 06:00 UTC daily Mon-Fri (before any market open)
0 6 * * 1-5 cd /app && python3 -m python_brain.ouroboros.universe_refresh >> /var/log/universe_refresh.log 2>&1

# Smart Ticker Selector — 06:30 UTC daily Mon-Fri (full daily re-score after universe refresh)
30 6 * * 1-5 cd /app && python3 -m python_brain.ouroboros.ticker_selector >> /var/log/ticker_selector.log 2>&1

# UNIFIED 15-minute ticker rotation — runs 22 hours/day (23:00-21:00 UTC Mon-Fri)
# Single unified watchlist across ALL 6 markets. Market-hours-aware:
# only includes tickers from exchanges currently open + LSE leveraged ETPs always.
# Ouroboros learns and improves rankings every 15 min with live price data.
# NO --skip-fetch: always uses live yfinance data for accurate scoring.

# Asian + pre-market hours: every 15 min from 23:00-07:45 UTC (Mon-Fri)
# 23:xx runs on Sunday night for Monday Asian open
*/15 23 * * 0-4 cd /app && python3 -m python_brain.ouroboros.ticker_selector >> /var/log/ticker_selector_15m.log 2>&1
*/15 0-7 * * 1-5 cd /app && python3 -m python_brain.ouroboros.ticker_selector >> /var/log/ticker_selector_15m.log 2>&1

# European + US overlap + US-only hours: every 15 min from 08:00-20:45 UTC (Mon-Fri)
*/15 8-20 * * 1-5 cd /app && python3 -m python_brain.ouroboros.ticker_selector >> /var/log/ticker_selector_15m.log 2>&1

# Session Briefing PDFs + Telegram (at session open)
# Asian session briefing: 00:55 UTC Mon-Fri (5 min before Asian scan starts)
55 0 * * 1-5 cd /app && python3 -m python_brain.ouroboros.session_pdf --session asian --send-telegram >> /var/log/session_pdf.log 2>&1

# European session briefing: 07:55 UTC Mon-Fri (5 min before European scan)
55 7 * * 1-5 cd /app && python3 -m python_brain.ouroboros.session_pdf --session european --send-telegram >> /var/log/session_pdf.log 2>&1

# American session briefing: 14:25 UTC Mon-Fri (5 min before US open)
25 14 * * 1-5 cd /app && python3 -m python_brain.ouroboros.session_pdf --session american --send-telegram >> /var/log/session_pdf.log 2>&1

# US-only session briefing: 16:30 UTC Mon-Fri (EU markets closed, US-only remainder)
30 16 * * 1-5 cd /app && python3 -m python_brain.ouroboros.session_pdf --session us_only --send-telegram >> /var/log/session_pdf.log 2>&1

# Daily Simulated Trade Report — 21:15 UTC Mon-Fri (just after DARK window starts)
15 21 * * 1-5 cd /app && python3 -m python_brain.ouroboros.daily_sim_report --send-telegram >> /var/log/sim_report.log 2>&1

# Telegram heartbeat — every 4 hours during weekdays
0 0,4,8,12,16,20 * * 1-5 cd /app && python3 -m python_brain.ouroboros.telegram_notify --heartbeat >> /var/log/telegram.log 2>&1

# FX Rate Refresh — every 6 hours Mon-Fri (live yfinance rates → fx_rates.toml)
0 0,6,12,18 * * 1-5 cd /app && python3 -m python_brain.ouroboros.fx_refresh >> /var/log/fx_refresh.log 2>&1

# Google Sheets Sync — every 5 min Mon-Fri (drain Redis queue → spreadsheet)
*/5 * * * 1-5 cd /app && python3 -m python_brain.ouroboros.sheets_sync >> /var/log/sheets_sync.log 2>&1

# Contract Expander — every 6 hours Mon-Fri (auto-grows contracts.toml from scored universe)
# Runs after ticker_selector has scored the universe. Validates via yfinance, appends new
# contracts, and sends SIGHUP to engine for hot-reload.
0 1,7,13,19 * * 1-5 cd /app && python3 -m python_brain.ouroboros.contract_expander >> /var/log/contract_expander.log 2>&1

# RT3: Cost Drag Report — 21:20 UTC Mon-Fri (after daily sim report, friction analysis)
20 21 * * 1-5 cd /app && python3 -m python_brain.ouroboros.cost_drag_report --send-telegram >> /var/log/cost_drag.log 2>&1

# RT2-P2: Fill Quality Report — 04:54 UTC Mon-Fri (nightly, paper→live slippage estimate)
54 4 * * 1-5 cd /app && python3 -m python_brain.ouroboros.fill_quality >> /var/log/fill_quality.log 2>&1

# N10w-N10gg: Post-Trade Diagnostics — 04:55 UTC Mon-Fri (nightly consolidated diagnostics)
55 4 * * 1-5 cd /app && python3 -m python_brain.ouroboros.post_trade_diagnostics >> /var/log/post_trade_diagnostics.log 2>&1

# IBKR Full Universe Scanner — 22:00 UTC every Sunday (weekly deep scan)
0 22 * * 0 cd /app && python3 -m python_brain.ouroboros.ibkr_scanner >> /var/log/ibkr_scanner.log 2>&1

# PHASE_8 #6: Ouroboros Monitor — 05:10 UTC Mon-Fri (20 min after nightly, checks TOML health + staleness)
10 5 * * 1-5 cd /app && python3 -m python_brain.ouroboros.ouroboros_monitor --all --send-telegram >> /var/log/ouroboros_monitor.log 2>&1

# v19-P2-15: Parquet Orphan Cleanup — 04:40 UTC daily (before nightly, prevents disk fill)
40 4 * * * cd /app && python3 -m python_brain.ouroboros.maintenance --cleanup >> /var/log/maintenance.log 2>&1

# P2-7: Write FTT Registry — 04:52 UTC Mon-Fri (after nightly_v6, before config_writer next day)
# Only needs to run once; config_fixes writes transaction_tax_registry.toml atomically
52 4 * * 1 cd /app && python3 -m python_brain.ouroboros.config_fixes --write-ftt >> /var/log/config_fixes.log 2>&1
