# nzt48-aegis-v2/crontab
# RT2: Bridge Health Monitor — every 15 min Mon-Fri during market hours (Telegram alert if unhealthy)
*/15 8-20 * * 1-5 cd /app && python3 -m python_brain.ouroboros.bridge_health >> /var/log/bridge_health.log 2>&1

# N10b: External Monitor — every 5 min Mon-Fri (deep health check: engine, disk, IBKR, Redis, cron)
*/5 * * * 1-5 cd /app && python3 -m python_brain.ouroboros.external_monitor --quiet >> /var/log/external_monitor.log 2>&1

# N10b: Daily Health Report — 08:00 UTC Mon-Fri (Telegram summary)
0 8 * * 1-5 cd /app && python3 -m python_brain.ouroboros.external_monitor --daily >> /var/log/external_monitor.log 2>&1

# N10c: Log Rotation — 04:45 UTC daily (before nightly learning, keeps 7 days)
45 4 * * * cd /app && python3 -m python_brain.ouroboros.log_rotate >> /var/log/log_rotate.log 2>&1

# Sprint S05 (H1): Sequential Nightly Pipeline — replaces individual entries
# Order: nightly_v6 → config_writer → win_loss_delta → claude_review → challenger → approval_gate
# flock prevents parallel execution. set -euo pipefail aborts on critical failures.
50 4 * * 1-5 flock -n /tmp/nightly.lock /app/scripts/nightly_pipeline.sh >> /var/log/nightly_pipeline.log 2>&1

# N6b: Claude Morning Briefing — 07:45 UTC Mon-Fri (before European session, uses N6a output)
45 7 * * 1-5 cd /app && python3 -m python_brain.ouroboros.claude_briefing --send-telegram >> /var/log/claude_briefing.log 2>&1

# Ouroboros v6.0 Backfill Simulator — 07:00 UTC daily Mon-Fri (pre-market learning)
0 7 * * 1-5 cd /app && python3 -m python_brain.ouroboros.backfill_simulator >> /var/log/ouroboros_backfill.log 2>&1

# ISA Universe Refresh — 06:00 UTC daily Mon-Fri (before any market open)
0 6 * * 1-5 cd /app && python3 -m python_brain.ouroboros.universe_refresh >> /var/log/universe_refresh.log 2>&1

# Universe Scanner — every 2 hours, 22h/day (23:00-21:00 UTC Mon-Fri)
# Provides top 100 tickers based on which markets are open.
# During LSE hours: leveraged/inverse ETPs boosted via configurable score bonus.
# flock prevents overlapping runs if previous scan hangs.
# Dark window = 21:00-23:00 UTC (2h, no scan).
0 23 * * 0-4 flock -n /tmp/ticker_selector.lock cd /app && python3 -m python_brain.ouroboros.ticker_selector >> /var/log/ticker_selector.log 2>&1
0 1,3,5,7,9,11,13,15,17,19,21 * * 1-5 flock -n /tmp/ticker_selector.lock cd /app && python3 -m python_brain.ouroboros.ticker_selector >> /var/log/ticker_selector.log 2>&1

# Bridge stderr log rotation — daily at 05:00 UTC (archive then truncate)
0 5 * * * cp /app/data/bridge_stderr.log /app/data/bridge_stderr.$(date +\%Y\%m\%d).log 2>/dev/null; truncate -s 0 /app/data/bridge_stderr.log 2>/dev/null; find /app/data -name 'bridge_stderr.2*.log' -mtime +7 -delete 2>/dev/null

# Session Briefing PDFs + Telegram (at session open)
# DST-aware: cron fires at BOTH possible UTC times (GMT and BST).
# dst_wrapper checks current DST state and only runs at the correct time.
#
# | Session   | London local | GMT (UTC) | BST (UTC) |
# |-----------|-------------|-----------|-----------|
# | asian     | 00:55       | 00:55     | 23:55     |
# | european  | 07:55       | 07:55     | 06:55     |
# | american  | 14:25       | 14:25     | 13:25     |
# | us_only   | 16:30       | 16:30     | 15:30     |

# Asian session briefing: 00:55 London (fires at 23:55 and 00:55 UTC)
55 23 * * 0-4 cd /app && python3 -m python_brain.ouroboros.dst_wrapper --session asian --send-telegram >> /var/log/session_pdf.log 2>&1
55 0 * * 1-5 cd /app && python3 -m python_brain.ouroboros.dst_wrapper --session asian --send-telegram >> /var/log/session_pdf.log 2>&1

# European session briefing: 07:55 London (fires at 06:55 and 07:55 UTC)
55 6 * * 1-5 cd /app && python3 -m python_brain.ouroboros.dst_wrapper --session european --send-telegram >> /var/log/session_pdf.log 2>&1
55 7 * * 1-5 cd /app && python3 -m python_brain.ouroboros.dst_wrapper --session european --send-telegram >> /var/log/session_pdf.log 2>&1

# American session briefing: 14:25 London (fires at 13:25 and 14:25 UTC)
25 13 * * 1-5 cd /app && python3 -m python_brain.ouroboros.dst_wrapper --session american --send-telegram >> /var/log/session_pdf.log 2>&1
25 14 * * 1-5 cd /app && python3 -m python_brain.ouroboros.dst_wrapper --session american --send-telegram >> /var/log/session_pdf.log 2>&1

# US-only session briefing: 16:30 London (fires at 15:30 and 16:30 UTC)
30 15 * * 1-5 cd /app && python3 -m python_brain.ouroboros.dst_wrapper --session us_only --send-telegram >> /var/log/session_pdf.log 2>&1
30 16 * * 1-5 cd /app && python3 -m python_brain.ouroboros.dst_wrapper --session us_only --send-telegram >> /var/log/session_pdf.log 2>&1

# Daily Simulated Trade Report — 21:15 UTC Mon-Fri (just after DARK window starts)
15 21 * * 1-5 cd /app && python3 -m python_brain.ouroboros.daily_sim_report --send-telegram >> /var/log/sim_report.log 2>&1

# Telegram heartbeat — every 4 hours during weekdays
0 0,4,8,12,16,20 * * 1-5 cd /app && python3 -m python_brain.ouroboros.telegram_notify --heartbeat >> /var/log/telegram.log 2>&1

# FX Rate Refresh — every 6 hours Mon-Fri (live yfinance rates → fx_rates.toml)
0 0,6,12,18 * * 1-5 cd /app && python3 -m python_brain.ouroboros.fx_refresh >> /var/log/fx_refresh.log 2>&1

# Google Sheets Sync — every 5 min Mon-Fri (drain Redis queue → spreadsheet)
*/5 * * * 1-5 cd /app && python3 -m python_brain.ouroboros.sheets_sync >> /var/log/sheets_sync.log 2>&1

# Contract Expander — every 6 hours Mon-Fri (auto-grows contracts.toml from scored universe)
# Runs after ticker_selector has scored the universe. Validates via yfinance, appends new
# contracts, and sends SIGHUP to engine for hot-reload.
0 1,7,13,19 * * 1-5 cd /app && python3 -m python_brain.ouroboros.contract_expander >> /var/log/contract_expander.log 2>&1

# RT3: Cost Drag Report — 21:20 UTC Mon-Fri (after daily sim report, friction analysis)
20 21 * * 1-5 cd /app && python3 -m python_brain.ouroboros.cost_drag_report --send-telegram >> /var/log/cost_drag.log 2>&1

# RT2-P2: Fill Quality Report — 04:54 UTC Mon-Fri (nightly, paper→live slippage estimate)
# flock: serialize behind nightly_pipeline to avoid WAL read race conditions
54 4 * * 1-5 flock -x /tmp/nightly.lock cd /app && python3 -m python_brain.ouroboros.fill_quality >> /var/log/fill_quality.log 2>&1

# N10w-N10gg: Post-Trade Diagnostics — 04:55 UTC Mon-Fri (nightly consolidated diagnostics)
# flock: serialize behind nightly_pipeline to avoid WAL read race conditions
55 4 * * 1-5 flock -x /tmp/nightly.lock cd /app && python3 -m python_brain.ouroboros.post_trade_diagnostics >> /var/log/post_trade_diagnostics.log 2>&1

# IBKR Full Universe Scanner — 22:00 UTC every Sunday (weekly deep scan)
0 22 * * 0 cd /app && python3 -m python_brain.ouroboros.ibkr_scanner >> /var/log/ibkr_scanner.log 2>&1

# PHASE_8 #6: Ouroboros Monitor — 05:10 UTC Mon-Fri (20 min after nightly, checks TOML health + staleness)
10 5 * * 1-5 cd /app && python3 -m python_brain.ouroboros.ouroboros_monitor --all --send-telegram >> /var/log/ouroboros_monitor.log 2>&1

# v19-P2-15: Parquet Orphan Cleanup — 04:40 UTC daily (before nightly, prevents disk fill)
40 4 * * * cd /app && python3 -m python_brain.ouroboros.maintenance --cleanup >> /var/log/maintenance.log 2>&1

# P2-7: Write FTT Registry — 04:52 UTC Mon-Fri (after nightly_v6, before config_writer next day)
# Only needs to run once; config_fixes writes transaction_tax_registry.toml atomically
52 4 * * 1 cd /app && python3 -m python_brain.ouroboros.config_fixes --write-ftt >> /var/log/config_fixes.log 2>&1

# === PLAN 2 CLAUDE INTELLIGENCE ROLES (Sprint S11-S19) ===

# S10: Claude Forensic Review — 04:53 UTC Mon-Fri (after nightly_v6, trade + veto classification)
# flock: serialize behind nightly_pipeline to avoid WAL read race conditions
53 4 * * 1-5 flock -x /tmp/nightly.lock cd /app && python3 -m python_brain.ouroboros.claude_forensic_review >> /var/log/claude_forensic.log 2>&1

# S12: Evening Briefing — 21:30 UTC Mon-Fri (after US close)
30 21 * * 1-5 cd /app && python3 -m python_brain.ouroboros.claude_briefing --evening --send-telegram >> /var/log/claude_briefing.log 2>&1

# S13: Universe Curation Shadow — every 2h during trading hours (shadow mode)
0 23 * * 0-4 cd /app && python3 -m python_brain.ouroboros.claude_curation >> /var/log/claude_curation.log 2>&1
0 1,3,5,7,9,11,13,15,17,19,21 * * 1-5 cd /app && python3 -m python_brain.ouroboros.claude_curation >> /var/log/claude_curation.log 2>&1

# S14: Weekly Gate Calibration — Friday 22:00 UTC
0 22 * * 5 cd /app && python3 -m python_brain.ouroboros.claude_rejected_review --send-telegram >> /var/log/claude_rejected_review.log 2>&1

# S18: Operator Psych Audit (H6) — Sunday 23:00 UTC
0 23 * * 0 cd /app && python3 -m python_brain.ouroboros.claude_psych_audit --send-telegram >> /var/log/claude_psych_audit.log 2>&1

# S19: SEC/RNS Filing Scanner (H7) — daily 06:00 UTC Mon-Fri (shadow mode)
0 6 * * 1-5 cd /app && python3 -m python_brain.ouroboros.claude_filing_scanner >> /var/log/claude_filing_scanner.log 2>&1

# === TIER 2: GEMINI PRO — Universe Curation (Sprint Intelligence) ===

# Gemini Core Universe Scan — every 2 hours Mon-Fri (aligned with ticker_selector, 5 min offset)
# Recommends top 80 session-aware tickers for streaming. Output: /app/data/gemini/core_universe_latest.json
5 23 * * 0-4 cd /app && GEMINI_API_KEY="${GEMINI_API_KEY}" python3 -m python_brain.ouroboros.gemini_scanner --core >> /app/data/gemini_scanner.log 2>&1
5 1,3,5,7,9,11,13,15,17,19,21 * * 1-5 cd /app && GEMINI_API_KEY="${GEMINI_API_KEY}" python3 -m python_brain.ouroboros.gemini_scanner --core >> /app/data/gemini_scanner.log 2>&1

# Gemini Morning Brief — 06:00 UTC Mon-Fri (before European session)
# Pre-market analysis: focus tickers, avoids, regime assessment, strategy weights.
0 6 * * 1-5 cd /app && GEMINI_API_KEY="${GEMINI_API_KEY}" python3 -m python_brain.ouroboros.gemini_scanner --brief >> /app/data/gemini_scanner.log 2>&1
